# deploy/gunicorn.service
[Unit]
Description=Cybersecurity Tools Web App
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/cybersecurity-tools-webapp
Environment="PATH=/var/www/cybersecurity-tools-webapp/.venv/bin"
ExecStart=/var/www/cybersecurity-tools-webapp/.venv/bin/gunicorn -c deploy/gunicorn.conf.py app.app:app
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
